# aio/rootfs/etc/s6-overlay/s6-rc.d/init-secrets/up
/etc/s6-overlay/scripts/00-init-secrets
