# aio/rootfs/etc/s6-overlay/s6-rc.d/sidekiq/run
#!/command/with-contenv bash
# Sidekiq service
# Depends on: postgresql, redis, rails

exec 2>&1

cd /app
DATA_DIR="${DELTABADGER_DATA_DIR:-/data}"

# Source secrets
if [ -f "${DATA_DIR}/secrets/.env.secrets" ]; then
    set -a
    source "${DATA_DIR}/secrets/.env.secrets"
    set +a
fi

# Brief delay to let Rails fully start
sleep 5

echo "Starting Sidekiq..."
exec bundle exec sidekiq -C config/sidekiq.yml
