# aio/rootfs/etc/s6-overlay/s6-rc.d/init-postgres/up
/etc/s6-overlay/scripts/01-init-postgres
