# aio/rootfs/etc/s6-overlay/s6-rc.d/rails/run
#!/command/with-contenv bash
# Rails (Puma) service
# Depends on: postgresql, redis, init-database

exec 2>&1

cd /app
DATA_DIR="${DELTABADGER_DATA_DIR:-/data}"

# Source secrets
if [ -f "${DATA_DIR}/secrets/.env.secrets" ]; then
    set -a
    source "${DATA_DIR}/secrets/.env.secrets"
    set +a
fi

# Clean up stale PID
rm -f tmp/pids/server.pid

# Symlink storage to data volume if not already done
if [ ! -L /app/storage ]; then
    rm -rf /app/storage 2>/dev/null || true
    ln -sf "$DATA_DIR/app/storage" /app/storage
fi

echo "Starting Puma..."
exec bundle exec puma -C config/puma.rb
